# packages/configure_strongswan/configure_strongswan/configure-strongswan
#!/bin/bash
# chkconfig: 2345 47 99
# processname: configure-strongswan
# description: setup strongswan configuration

. /etc/rc.d/init.d/functions

SUBJ_PROG=/usr/bin/gen_subj.sh
PASSWD_PROG=/usr/bin/gen_passwd.sh

WORKDIR="/tmp/strongswan"
CA_KEY="$WORKDIR/rootCA.key"
CA_PEM="$WORKDIR/ca.pem"
CA_P12="$WORKDIR/ca.crt"
SERVER_KEY="$WORKDIR/server.key"
SERVER_PEM="$WORKDIR/server.pem"
IPSEC_CONF="$WORKDIR/ipsec.conf"
IPSEC_SECRETS="$WORKDIR/ipsec.secrets"
IPSEC_CONF_TMPL="$WORKDIR/ipsec.conf.tmpl"
IPSEC_SECRETS_TMPL="$WORKDIR/ipsec.secrets.tmpl"
CA_SUBJ=`$SUBJ_PROG`
CA_PASSWD=`$PASSWD_PROG`
SERVER_SUBJ=`$SUBJ_PROG`
SERVER_PASSWD=`$PASSWD_PROG`
#TODO allow the user to specify the subnet to use - definitely need this
RIGHT_IP="192.168.31.0\/24"
STRONGSWANDIR=/etc/strongswan
IPSECDIR=$STRONGSWANDIR/ipsec.d
CERTDIR=$IPSECDIR/certs
PRIVATEDIR=$IPSECDIR/private
CA_PASSWD_FILE=$PRIVATEDIR/ca_pass.txt
CACERTDIR=$IPSECDIR/cacerts
DAYS=1024
KEYSIZE=4096

install_files() {
	cp $CA_PEM $CACERTDIR
	cp $SERVER_PEM $CERTDIR
	cp $SERVER_KEY $PRIVATEDIR
	cp $CA_KEY $PRIVATEDIR
	cp $IPSEC_CONF $STRONGSWANDIR
	cp $IPSEC_SECRETS $STRONGSWANDIR
	cp $CA_P12 $CACERTDIR
	echo $CA_PASSWD > $CA_PASSWD_FILE
}

# args $1 - filename, $2 password
gen_key() {
	openssl genrsa -passout pass:$2 -out $1 $KEYSIZE
}

# args $1 - filename, $2 keyfile, $3 subject
self_sign() {
	openssl req -x509 -new -nodes -key $2 -days $DAYS -out $1 -subj "$3"
	openssl x509 -inform PEM -outform DER -in $1 -out $CA_P12
}

# args:
# $1 - key
# $2 - key password
# $3 - subject
# $4 - output file
ca_sign() {
	openssl req -new -key $1 -out request.csr -subj "$3" -passin pass:$2
	openssl x509 -req -in request.csr -CA $CA_PEM -CAkey $CA_KEY -CAcreateserial -out $4 -days $DAYS -passin pass:$2
}

# args
# $1 server key pass
update_ipsec_secrets() {
	sed "s/{SERVERPASS}/$1/" $IPSEC_SECRETS_TMPL > $IPSEC_SECRETS
}

# args
# $1 server IP
# $2 tun IP address assigned to client
update_ipsec_conf() {
        sed "s/{LEFTIP}/$1/" $IPSEC_CONF_TMPL > $IPSEC_CONF
	sed -i "s/{RIGHTIP}/$2/" $IPSEC_CONF
}

start() {
  mkdir -p $WORKDIR
  cd $WORKDIR

  cp /usr/share/strongswan/templates/config/ipsec.conf.tmpl $WORKDIR/
  cp /usr/share/strongswan/templates/config/ipsec.secrets.tmpl $WORKDIR/ 
  gen_key $CA_KEY $CA_PASSWD
  self_sign $CA_PEM $CA_KEY "$CA_SUBJ"

  gen_key $SERVER_KEY $SERVER_PASSWD
  ca_sign $SERVER_KEY $SERVER_PASSWD "$SERVER_SUBJ" $SERVER_PEM

  #we run after networking is up
  SERVER_IP=`ifconfig eth0 | grep "inet addr" | cut -d: -f2 | cut -d" " -f1`

  update_ipsec_conf $SERVER_IP $RIGHT_IP
  update_ipsec_secrets $SERVER_PASSWD

  install_files

  rm -rf $WORKDIR
}

stop() {
  echo "Nothing to do here"
}
 
restart() {
  stop
  start
}

# See how we were called.
case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    restart
    ;;
  *)
    echo $"Usage: $0 {start|stop|restart}"
    exit 1
esac
 
exit $?

